# modules/systemd/dotfiles/systemd/user/dynamic-window-rules.service
[Unit]
Description=Dynamic Window Rules
BindsTo=niri.service
After=niri.service

[Service]
ExecStart=/home/jappe/.local/bin/dynamic_window_rules.py
